# 10-hci.rules
ACTION=="add", KERNEL=="hci*", RUN+="/usr/sbin/sixad-helper"

